# CSCI 330/a7partb.cron
0 5 * * 1 checkgrades
0 5 * * 5 checkgrades
